# debian/freeradius.prerm
#! /bin/sh

set -e

case "$1" in
  remove)
	;;
  upgrade)
	;;
  failed-upgrade)
	;;
  deconfigure)
	;;
esac

if [ "$1" = "upgrade" -o  "$1" = "remove" ] ; then
	if [ -x /usr/sbin/invoke-rc.d ] ; then
		invoke-rc.d freeradius stop
	else
		/etc/init.d/freeradius stop
	fi
fi

exit 0
