# debian/freeradius-postgresql.postinst
#! /bin/sh

set -e

case "$1" in
  configure)
	if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
		invoke-rc.d freeradius restart
	else
		/etc/init.d/freeradius restart
	fi
	;;
  abort-upgrade)
	;;
  abort-remove)
	;;
  abort-deconfigure)
	;;
esac

#DEBHELPER#
