# debian/radiusd-freeradius.conffiles
etc/cron.daily/freeradius
etc/cron.monthly/freeradius
etc/init.d/freeradius
